# spdynu.service
[Unit]
Description=spdyn Updater

[Service]
Type=simple
ExecStart=/bin/spdynu

[Install]
WantedBy=multi-user.target
