# aur/spdynu.service
[Unit]
Description=spdyn Updater

[Service]
Type=simple
ExecStart=/usr/bin/spdynu

[Install]
WantedBy=multi-user.target
